# deploy/supervisor/eventche-api.conf
; EvenTche API — Supervisor
; Copie para /etc/supervisor/conf.d/eventche-api.conf
; Depois: supervisorctl reread && supervisorctl update && supervisorctl start eventche-api:*
;
; Pré-requisito (deploy / boot, uma vez):
;   cd /var/www/api && php artisan rabbitmq:setup-queues

[group:eventche-api]
programs=api-scheduler,api-worker-default,api-worker-payments-create,api-worker-payments-webhook,api-worker-emails,api-worker-tickets-generation,api-worker-tickets-expiration
priority=999

; ---------------------------------------------------------------------------
; Scheduler (Laravel 12) — executa comandos agendados em routes/console.php
; ---------------------------------------------------------------------------
[program:api-scheduler]
process_name=%(program_name)s
command=php /var/www/api/artisan schedule:work
directory=/var/www/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/scheduler.log
stopwaitsecs=60

; ---------------------------------------------------------------------------
; Workers RabbitMQ — uma fila por programa (QUEUE_CONNECTION=rabbitmq)
; --max-time=3600 recicla o processo a cada hora (evita vazamento de memória)
; ---------------------------------------------------------------------------

[program:api-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan queue:work rabbitmq --queue=default --sleep=3 --tries=3 --timeout=90 --max-time=3600
directory=/var/www/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/worker-default.log
stopwaitsecs=3600

[program:api-worker-payments-create]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan queue:work rabbitmq --queue=payments.create --sleep=3 --tries=5 --timeout=120 --max-time=3600
directory=/var/www/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/worker-payments-create.log
stopwaitsecs=3600

[program:api-worker-payments-webhook]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan queue:work rabbitmq --queue=payments.webhook --sleep=1 --tries=5 --timeout=120 --max-time=3600
directory=/var/www/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/worker-payments-webhook.log
stopwaitsecs=3600

[program:api-worker-emails]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan queue:work rabbitmq --queue=emails --sleep=3 --tries=3 --timeout=90 --max-time=3600
directory=/var/www/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/worker-emails.log
stopwaitsecs=3600

[program:api-worker-tickets-generation]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan queue:work rabbitmq --queue=tickets.generation --sleep=3 --tries=3 --timeout=180 --max-time=3600
directory=/var/www/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/worker-tickets-generation.log
stopwaitsecs=3600

[program:api-worker-tickets-expiration]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan queue:work rabbitmq --queue=tickets.expiration --sleep=3 --tries=3 --timeout=90 --max-time=3600
directory=/var/www/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/worker-tickets-expiration.log
stopwaitsecs=3600
